# deploy/systemd/sim-display.service
# Headless X server for the sim — NVIDIA renders without a physical monitor.
# Requires: /etc/X11/xorg.conf generated by
#   sudo nvidia-xconfig --allow-empty-initial-configuration
# Install: sudo cp sim-display.service /etc/systemd/system/ && sudo systemctl enable --now sim-display
[Unit]
Description=Headless X display :1 for DroneSim
After=multi-user.target

[Service]
ExecStart=/usr/bin/X :1 -nolisten tcp vt7
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
